# systemd/misc-x@.service
[Unit]
Description=Miscellaneous settings for X11 on %I
PartOf=graphical-session@%i.target

[Service]
Environment=DISPLAY=%I
ExecStart=/usr/bin/xset -b
ExecStart=%h/.config/awesome/bin/xdg-defaults
Type=oneshot
RemainAfterExit=false
