# systemd/wallpaper@.timer
[Unit]
Description=Wallpaper rotation for %I
PartOf=graphical-session@%i.target

[Timer]
OnCalendar=hourly
RandomizedDelaySec=10m
