# systemd/xiccd@.service
[Unit]
Description=X color management for %I
PartOf=graphical-session@%i.target

[Service]
Environment=DISPLAY=%I
ExecStart=/usr/bin/xiccd --edid
