# systemd/xsession@.target
[Unit]
Description=X session on %I
BindsTo=graphical-session@%i.target
