# systemd/xsession@.target
[Unit]
Description=X session managed by systemd
BindsTo=graphical-session@%i.target
